# scripts/yai.service
[Unit]
Description=yAI Application Service
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/yai
Environment="PATH=/home/ubuntu/yai/venv/bin:/usr/local/bin:/usr/bin:/bin:/sbin:/usr/sbin:/usr/local/sbin"
Environment="LD_LIBRARY_PATH=/home/ubuntu/yai/boost_install/lib"
Environment="YAI_CHAT_INPUT_PATH=sinput"
ExecStart=gunicorn --workers 2 --bind 0.0.0.0:80 yai.wsgi:application
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
